# update-localrepo.timer
[Unit]
Description=Run local apt repo updater weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
